# upstart/occi.conf

start on local-filesystems

task

script
	exec /usr/bin/occi
end script
