# upstart/rpi_halt_button.conf
start on (local-filesystems and net-device-up IFACE=lo)
stop on [!12345]
script
    python /home/pi/src/LittleSleeper/rpi-halt-btn.py
end script

